# Docker/root/etc/services.d/chinesesubfinder/run
#!/usr/bin/with-contenv bash

# 启动chinesesubfinder
cd /app
exec s6-setuidgid abc chinesesubfinder
