# docker/full-rootfs/etc/services.d/xvfb/run
#!/usr/bin/with-contenv bash

exec Xvfb -ac ${DISPLAY} -screen 0 1280x1024x16
