# docker/normal-rootfs/etc/services.d/chinesesubfinder/run
#!/usr/bin/with-contenv bash

umask ${UMASK}

cd /config 
exec s6-setuidgid ${PUID}:${PGID} chinesesubfinder
